# testdata/container-with-service/expected/app.container

[Container]
Image=docker.io/myapp:v2
ContainerName=app
Environment=APP_ENV=production
Environment=LOG_LEVEL=info
Volume=app-data.volume:/data
DropCapability=ALL
NoNewPrivileges=true
Memory=512m
HealthCmd=curl -f http://localhost/health
HealthInterval=30s
HealthRetries=3
HealthTimeout=5s

[Service]
Restart=always
CPUQuota=50%
MemoryMax=1G

[Install]
WantedBy=multi-user.target
