# testdata/pod-with-containers/expected/webapp-sidecar.container

[Container]
Image=docker.io/envoy:latest
ContainerName=webapp-sidecar
Pod=webapp.pod

[Install]
WantedBy=multi-user.target
