# testdata/pod-with-containers/expected/webapp.container

[Container]
Image=docker.io/myapp:latest
ContainerName=webapp-app
Pod=webapp.pod

[Install]
WantedBy=multi-user.target
