# testdata/cross-quadlet-ref/expected/proxy.container
[Unit]
Requires=app.service
After=app.service

[Container]
Image=docker.io/haproxy:latest
ContainerName=proxy
Network=proxy-backend.network

[Install]
WantedBy=multi-user.target
