# testdata/json-label/expected/app.container

[Container]
Image=docker.io/app:latest
Label=plain=1
Label=app.spec={"kind":"container"}
