# testdata/resource-name-ref/expected/cache.container

[Container]
Image=docker.io/redis:7

[Install]
WantedBy=default.target
