# testdata/resource-name-ref/expected/db.container

[Container]
Image=docker.io/postgres:16
ContainerName=postgres

[Install]
WantedBy=default.target
