# testing/quadlet-section/expected/nodeps.container

[Container]
Image=docker.io/standalone:latest

[Quadlet]
DefaultDependencies=false

[Install]
WantedBy=default.target
