# testdata/unit-install-full/expected/fullunit.container
[Unit]
Description=Full unit test container
Documentation=https://example.com/docs
SourcePath=/etc/systemd/system/fullunit.container
Requires=redis.service
Wants=network-online.target
BindsTo=redis.service
Conflicts=legacy-app.service
Before=app-ready.target
After=network-online.target
After=redis.service
OnSuccess=notify-success.service
OnFailure=notify-failure.service
PartOf=app.target
StopWhenUnneeded=true
RefuseManualStart=true
RefuseManualStop=true
AllowIsolate=true
OnSuccessJobMode=fail
OnFailureJobMode=replace
IgnoreOnIsolate=true
StartLimitIntervalSec=60s
StartLimitBurst=5
ConditionPathExists=/etc/app.conf
AssertPathExists=/opt/app/bin

[Container]
Image=docker.io/fullunit:latest

[Install]
Alias=myapp.service
WantedBy=multi-user.target
RequiredBy=app.target
UpheldBy=core.target
DefaultInstance=default
