# testdata/basic-container/expected/nginx.container

[Container]
Image=docker.io/nginx:latest
ContainerName=nginx

[Install]
WantedBy=default.target
