# testdata/checks/expected/app.container

[Container]
Image=docker.io/app:latest
PublishPort=8080:80
