# testdata/service-section-full/expected/svcfull.container

[Container]
Image=docker.io/app:latest

[Service]
Type=notify
Restart=on-failure
RestartSec=5s
TimeoutStartSec=30s
TimeoutStopSec=15s
ExecStartPre=/usr/bin/check-prereqs
ExecStartPost=/usr/bin/notify-ready
ExecStop=/usr/bin/graceful-stop
ExecStopPost=/usr/bin/cleanup
ExecReload=/usr/bin/reload-config
WatchdogSec=60s
RemainAfterExit=true
KillMode=mixed
KillSignal=SIGTERM
NotifyAccess=main
RestartPreventExitStatus=1 6 SIGABRT
SuccessExitStatus=0 2
Environment=SVC_ENV=production
EnvironmentFile=/etc/svcfull.env
WorkingDirectory=/opt/svcfull
MemoryMax=2G
MemoryHigh=1G
MemoryLow=256M
MemoryMin=128M
CPUQuota=200%
CPUWeight=100
IOWeight=500
TasksMax=4096
LimitNOFILE=65536:131072
LimitNPROC=4096
LimitCORE=infinity
LimitMEMLOCK=8G

[Install]
WantedBy=multi-user.target
