# testdata/service-section-full/expected/svcfull.container

[Container]
Image=docker.io/app:latest

[Service]
ExecStartPre=/usr/bin/check-prereqs
ExecReload=/usr/bin/reload-config
ExecStartPost=/usr/bin/notify-ready
ExecStop=/usr/bin/graceful-stop
ExecStopPost=/usr/bin/cleanup
RestartSec=5s
TimeoutStartSec=30s
TimeoutStopSec=15s
WatchdogSec=60s
Type=notify
Restart=on-failure
RemainAfterExit=true
RestartPreventExitStatus=1
RestartPreventExitStatus=6
RestartPreventExitStatus=SIGABRT
SuccessExitStatus=0
SuccessExitStatus=2
NotifyAccess=main
WorkingDirectory=/opt/svcfull
Environment=SVC_ENV=production
EnvironmentFile=/etc/svcfull.env
LimitCORE=infinity
LimitNOFILE=65536:131072
LimitNPROC=4096
LimitMEMLOCK=8G
CPUWeight=100
CPUQuota=200%
MemoryMin=128M
MemoryLow=256M
MemoryHigh=1G
MemoryMax=2G
IOWeight=500
TasksMax=4096
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
