# testing/build-inline-containerfile-additional/expected/myapp.container

[Container]
Image=localhost/myapp:latest
ContainerName=myapp

[Install]
WantedBy=default.target
